# base/common/etc/uci-defaults/luci-theme-cortina
uci batch <<-EOF
        set luci.themes.Cortina=/luci-static/cortina
        commit luci
EOF
