# base.bak/common/etc/uci-defaults/luci-app-ar988x
#!/bin/sh
sed -i 's/has_wifi = true/has_wifi = true/g' /usr/lib/lua/luci/controller/admin/network.lua
