# base/base-0.1.0-ftm-100s/etc/uci-defaults/luci-app-dhcp-simple
#!/bin/sh
# Nothing to do (dhcp is already in ucitrack)
exit 0
